# resources/device_nanny.rules
#UDEV rule for /etc/udev/rules.d/
SUBSYSTEM=="usb" ENV{DEVTYPE}=="usb_device" RUN+="/home/pi/devicenannyv2/start_checkout.sh"